# chart_world_import.cron
#!/usr/bin/env bash
#
# run importer every 5 min
#
*/5 * * * * /usr/local/bin/chart_world_import > /dev/null 2>&1 
